# lightdm-failsafe/lib/systemd/system/lightdm-failsafe.service
[Unit]
Description=Use if lightdm fail
    
[Service]
Type=oneshot
ExecStart=/usr/sbin/lightdm-failsafe
